# valheim-monitor.service
[Service]
WorkingDirectory={{ RESOURCES }}
ExecStart=/bin/bash monitor.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
